# deploy/yt-transcript.service
[Unit]
Description=YT Transcript API
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/.openclaw/factory/worktrees/yt-transcript-web/backend
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8001
EnvironmentFile=/home/azureuser/.openclaw/factory/worktrees/yt-transcript-web/backend/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
